# etc/profile.d/swiftbot-paths.sh
#!/bin/bash
export LD_LIBRARY_PATH=$LD_LIBRARY_PATH:/opt/swiftbot/lib
export PATH=$PATH:/opt/swiftbot/bin
